# data/minecraftd.service.in
[Unit]
Description=Minecraft server daemon

[Service]
ExecStart=@bindir@/minecraftd
ExecStop=dbus-send --system --dest=net.za.slyfox.Minecraftd1 /net/za/slyfox/Minecraftd1 net.za.slyfox.Minecraftd1.Stop
BusName=net.za.slyfox.Minecraftd1

