# data/minecraftd.service.in
[Unit]
Description=Minecraft daemon

[Service]
ExecStart=@bindir@/minecraftd
BusName=net.za.slyfox.Minecraftd1

